# server/debian/moneydb-backup.service
[Unit]
Description=Backup MoneyDB

[Service]
Type=oneshot
ExecStart=/bin/bash /var/backups/backup-moneydb.sh
