# server/debian/moneydb-backup.timer
[Unit]
Description=Daily backup of MoneyDB

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target



